# praktikum-modul-1-b14/task-3/crontabs
*/2 * * * * /bin/bash /home/marcomanalu/Documents/praktikum-modul-1-b14/task-3/automation.sh
*/10 * * * * /bin/bash /home/marcomanalu/Documents/praktikum-modul-1-b14/task-3/download.sh
0 */2 * * * /bin/bash /home/marcomanalu/Documents/praktikum-modul-1-b14/task-3/archive.sh
